# USBIP/90-usb-policy.rules
SUBSYSTEM=="usb", \
  ACTION=="add", \
  ENV{DEVTYPE}=="usb_device", \
  RUN+="/usr/local/bin/usb-policy-handler.sh %k"
